# www-apps/pterodactyl-panel-common/files/pteroq
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

PANEL_INST=${SVCNAME#*.}

name="Pterodactyl Panel Queue Worker (${PANEL_INST})"
description=""
command=/usr/bin/php
command_args="${pteroq_panel_install}/artisan queue:work --queue=${pteroq_queue_order} --sleep=${pteroq_sleep} --tries=${pteroq_tries} ${pteroq_args}"
command_user="${pteroq_user}"
command_background=1
pidfile="/run/${SVCNAME}.pid"
output_log="/var/log/${SVCNAME}.log"
error_log="${output_log}"

depend() {
  need redis
}

start_pre() {
  if [ "${SVCNAME}" = "pteroq" ]; then
    eerror "You are not supposed to run this script directly. Create a symlink"
    eerror "for the Pterodactyl Panel you want to start the queue for as well as"
    eerror "a copy of the configuration file and modify it appropriately like so..."
    eerror
    eerror "  ln -s /etc/init.d/pteroq pteroq.myinstance"
    eerror "  cp /etc/conf.d/pteroq /etc/conf.d/pteroq.myinstance"
    eerror "  `basename "${EDITOR}"` /etc/conf.d/pteroq.myinstance"
    eerror
    return 1
  fi

  ebegin "Correcting ownership on log file"
  touch "${output_log}" "${error_log}" && \
    chown "${command_user}" "${output_log}" "${error_log}"
  eend $?
}
