# www-apps/srht-meta-api/files/srht-meta-api
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="meta.sr.ht API service"
description=""
command=/usr/bin/metasrht-api
command_args="-b ${SRHT_META_API_ADDR:-127.0.0.1:5100}"
command_background=1
pidfile="/var/run/${RC_SVCNAME}.pid"
command_user="srht-meta"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="${output_log}"

depend() {
  need net
  use redis postgresql
}
