# net-misc/pterodactyl-daemon/files/pterodactyl-daemon
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="pterodactyl daemon"
description=""
command=/usr/bin/npm
command_args="start"
command_user=pterodactyl
directory="/usr/lib/pterodactyl-daemon"
command_background=1
pidfile="/run/${SVCNAME}.pid"
output_log="/var/log/pterodactyl-daemon/output.log"
error_log="${output_log}"

depend() {
  need net docker
}
