# net-misc/pterodactyl-daemon/files/pterodactyl-daemon-1.0.1
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="pterodactyl daemon"
description=""
command=/usr/bin/wings
command_args="--config /etc/pterodactyl-daemon/config.yml"
command_user=pterodactyl
directory="/etc/pterodactyl-daemon"
command_background=1
pidfile="/run/${SVCNAME}.pid"
output_log="/var/log/pterodactyl-daemon.log"
error_log="${output_log}"

depend() {
  need net docker
}
