# www-apps/srht-paste/files/srht-paste-r1
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="paste.sr.ht website service"
description=""
command=/usr/bin/gunicorn
command_args="pastesrht.app:app -b ${SRHT_PASTE_ADDR:-127.0.0.1:5011}"
command_background=1
pidfile="/var/run/${RC_SVCNAME}.pid"
command_user="srht-paste"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="${output_log}"

depend() {
  need net
  use postgresql
}
