# www-apps/srht-paste/files/srht-paste-webhooks
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="paste.sr.ht webhooks service"
description=""
command=/usr/bin/celery
command_args="-A pastesrht.webhooks worker --loglevel=INFO"
command_user="srht-paste"
output_log="/var/log/srht-paste-webhooks.log"
error_log="${output_log}"

depend() {
  need net
  use postgresql
}
